# templates/phpdev-supervisor.in
#!/bin/sh
# RedHat startup script for a supervisor instance
#
# chkconfig: 2345 80 20
# description: ${:script-name} 

# Source function library.
. /etc/rc.d/init.d/functions

supervisorctl="${buildout:bin-directory}/supervisorctl"
supervisord="${buildout:bin-directory}/supervisord"
name="${:script-name}"

[ -f $supervisord ] || exit 1
[ -f $supervisorctl ] || exit 1

RETVAL=0

start() {
     echo -n "Starting $name: "
     $supervisord
     RETVAL=$?
     [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$name
     echo
     return $RETVAL
}

stop() {
     echo -n "Stopping $name: "
     $supervisorctl shutdown
     RETVAL=$?
     [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$name
     echo
     return $RETVAL
}

status() {
     echo "$name status:"
     $supervisorctl status
     RETVAL=$?
     [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$name
     echo
     return $RETVAL
}

case "$1" in
         start)
             start
             ;;

         stop)
             stop
             ;;

         status)
             status
             ;;

         restart)
             stop
             start
             ;;
esac

exit $REVAL
